# colana-ca1-bot.service
[Unit]
Description=Telegram Token Analysis Bot - Simple CA1 Version
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/projects/colana
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/root/projects/colana
EnvironmentFile=/root/projects/colana/config/.env
ExecStart=/usr/bin/python3 /root/projects/colana/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=colana-bot

# 安全设置
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/root/projects/colana/storage

[Install]
WantedBy=multi-user.target
