# deployment/colana-bot.service
[Unit]
Description=Solana Token Analysis Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/colana
ExecStart=/root/projects/colana/scripts/start.sh
Restart=always
RestartSec=10
Environment=PATH=/root/projects/colana/venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/root/projects/colana/config/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
